# deploy/systemd/edge-agent.service
[Unit]
Description=EdgeAI Edge Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/edge-agent/edge-agent.env
ExecStart=/usr/local/bin/edge-agent --config /etc/edge-agent/config.json
Restart=always
RestartSec=5
LimitNOFILE=65536
WorkingDirectory=/etc/edge-agent

[Install]
WantedBy=multi-user.target
